# dev_env/crontab
* * * * * /usr/local/bin/php /code/artisan schedule:run >> /dev/null 2>&1
